# Practices/hxpCTF2020-kernel-rop/initramfs/etc/inittab
::sysinit:/etc/init.d/rcS
::once:-sh -c 'cat /etc/motd; setuidgid 0 sh; poweroff'
